# ansible/templates/quicknotes-ansible-pull.timer.j2
[Unit]
Description=Run QuickNotes ansible-pull every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=quicknotes-ansible-pull.service

[Install]
WantedBy=timers.target
